# templates/default/xtreemfs-osd-all.conf.erb
description "XtreemFS OSD services"

setuid <%= @user %>
setgid <%= @group %>

start on <%= @start_on %>
stop on <%= @stop_on %>
